# deploy/systemd/fortune-bot-submit-window.timer
[Unit]
Description=Run Fortune Bot submit-window every 15 minutes

[Timer]
OnCalendar=*:0/15
Persistent=true
Unit=fortune-bot-submit-window.service

[Install]
WantedBy=timers.target
